# shellspecter.service
[Unit]
Description=ShellSpecter System Observatory Daemon
Documentation=https://github.com/ExythAI/ShellSpecter
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/opt/shellspecter/ShellSpecter.Specter
WorkingDirectory=/opt/shellspecter
User=specter
Group=specter
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Security hardening
ProtectHome=yes
ProtectSystem=full
ReadOnlyPaths=/proc /sys
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment=DOTNET_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5050

[Install]
WantedBy=multi-user.target
